# .config/autostart/steam_silent.desktop
[Desktop Entry]
Type=Application
Name=steam-silent
Exec=/usr/bin/steam-runtime -silent %U
Terminal=false
